# deploy/superego-domain-ready.timer
[Unit]
Description=Check Superego domain DNS every five minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
